# ceph-n1/.provision/sharedfolders.conf
start on vagrant-mounted
pre-start script
  /bin/mount -t vboxsf -o uid=10000,gid=10000 home_fengzhiguo_src /home/fengzhiguo/src
end script
